# systemd/motionflow.service
[Unit]
Description=MotionFlow - Pose-based action detection engine
After=network.target
Wants=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__

# Source the SDK venv activate script (sets AXELERA_FRAMEWORK, AXELERA_RUNTIME_DIR, etc.)
# Using 'exec' so python3 replaces bash as the main process (gets signals correctly)
ExecStart=/bin/bash -c 'source ~/voyager-sdk/venv/bin/activate && exec python3 main.py --config config/settings.yaml'

# Axelera SDK board configuration
Environment=AXELERA_CONFIGURE_BOARD=,30


# Exit code 42 = restart requested from WebUI (config change requiring full restart)
# SystemD will treat exit code 42 as a success but still restart due to RestartForceExitStatus
SuccessExitStatus=42
RestartForceExitStatus=42
RestartSec=2

# Restart on failure (crash, signal) and on exit code 42
Restart=on-failure

# Limit restart attempts to avoid tight crash loops
StartLimitIntervalSec=60
StartLimitBurst=5

# Resource limits
Nice=0
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=motionflow

[Install]
WantedBy=multi-user.target
