# systemd/motionflow-webui.service
[Unit]
Description=MotionFlow Web UI - Configuration & live preview
After=network.target
Wants=motionflow.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__

# Source the SDK venv activate script (sets PYTHONPATH and SDK env vars)
ExecStart=/bin/bash -c 'source ~/voyager-sdk/venv/bin/activate && exec python3 web_ui/app.py --config config/settings.yaml --port 5000'

# PYTHONPATH needed so web_ui/app.py can import from project root
Environment=PYTHONPATH=__WORKDIR__

# Restart on failure only
Restart=on-failure
RestartSec=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=motionflow-webui

[Install]
WantedBy=multi-user.target
